# layers/meta-balena-5x-owa/recipes-core/system_maintenance/files/system_maintenance.service
[Unit]
Description=Sytem maintenance service
After=pmsrv.service

[Service]
ExecStart=/usr/bin/system_maintenance.sh

[Install]
WantedBy=multi-user.target
