# layers/meta-balena-5x-owa/recipes-core/pmsrv/files/pmsrv.service
[Unit]
Description=pmsrv service
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/pmsrv
Restart=always

[Install]
WantedBy=switchwifi.service
